# config/ansible/files/rathole.service
[Unit]
Description=rathole - A lightweight reverse proxy for NAT traversal
After=network.target
Wants=network.target

[Service]
Type=exec
ExecStart=/usr/local/bin/rathole-client-start
Restart=always
RestartSec=5
User=root
Group=root

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /etc/rathole

[Install]
WantedBy=multi-user.target
